# dot_config/systemd/user/protonmail-bridge.service.d/override.conf.tmpl
[Unit]
{{- if (index . "yubikeyring" | default false) }}
Requires=yubikeyring.service
After=yubikeyring.service
{{ end }}
